# package/peridio/peridiod/peridiod.service
[Unit]
Description=Peridio Daemon

[Service]
Restart=always
ExecStart=/usr/sbin/peridiod start
User=root

[Install]
WantedBy=multi-user.target
